# supervisor.conf
[program:lsc]
command=/var/www/lishengchun/gunicorn_start
user=root
autostart=true
autorestart=true